# systemd/docker-compose@cn-roms.service
[Unit]
Description=cn-roms docker compose stack (RomM)
# Hard ordering: don't even try until the NFS share is mounted + docker is up.
Requires=docker.service home-gonzalo-docker-data-nfs.mount
# Serial start ordering with sibling stacks on this host: cn-bittorrent ->
# cn-media -> cn-roms. Three docker-compose@ units all touch the docker
# socket and the in-stack watchdog runs `docker compose up` concurrently;
# parallel-firing after boot races the watchdog/compose pairs into name
# conflicts. Sequential After= ordering eliminates the parallel trigger
# (it's ordering only, not a dependency — cn-roms still starts if
# cn-media fails).
After=docker.service home-gonzalo-docker-data-nfs.mount network-online.target docker-compose@cn-media.service
RequiresMountsFor=/home/gonzalo/docker/data/nfs
# BindsTo: if the NFS mount unit FAILS (e.g. mount timeout, as on 2026-06-22),
# tear the stack down instead of leaving RomM bound to the empty local
# placeholder (which leaves it with an empty library). `hard,_netdev` means a
# transient server blip doesn't change the .mount state, so this won't flap.
BindsTo=home-gonzalo-docker-data-nfs.mount
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=gonzalo
Group=gonzalo
WorkingDirectory=/home/gonzalo/cn-roms
# `${PWD}` is referenced by the watchdog's bind mounts (CLAUDE.md §15
# invariant: mount the project at the same absolute path on host and
# container). systemd does NOT set PWD automatically, so we have to.
Environment="PWD=/home/gonzalo/cn-roms"
# Belt-and-suspenders pre-checks: bail BEFORE docker compose ever sees
# the bind sources if the mount or the roms/ subtree is missing.
ExecStartPre=/usr/bin/mountpoint -q /home/gonzalo/docker/data/nfs
# Assert it's REALLY nfs, not the empty local placeholder a failed mount leaves.
ExecStartPre=/bin/sh -c '[ "$(stat -f -c %%T /home/gonzalo/docker/data/nfs)" = nfs ]'
ExecStartPre=/usr/bin/test -d /home/gonzalo/docker/data/nfs/roms/library/roms
ExecStart=/usr/bin/docker compose -p cn-roms up -d --remove-orphans
ExecStop=/usr/bin/docker compose -p cn-roms down

[Install]
WantedBy=multi-user.target
